# dev/crontab
*/5 * * * *     /bin/bash /app/dev/run_jobs.sh > /var/log/cron.log 2>&1
